# crontab
#!/bin/sh

sleep $(( RANDOM % 7200 )); /usr/local/bin/registerOnion.sh

